# nginx/conf/etc/supervisor/conf.d/dnsmasq.conf
[program:dnsmasq]
command=/usr/sbin/dnsmasq -d
autorestart=true
startretries=3