# apache2/conf/etc/supervisor/conf.d/nginx.conf
[program:nginx]
command=/usr/local/sbin/nginx -g "daemon off;"
autorestart=true
startretries=3